# systemd/sonarr.service.d/override.conf
[Unit]
Requires=mergerfs.service
After=network.target mergerfs.service
