# systemd/appdatabackup.service
[Unit]
Description=Appdata Backup Service

[Service]
User=ruxn
Group=ruxn
ExecStart=/opt/bin/appdatabackup.sh
Restart=no
