# systemd/cleartorrent.service
[Unit]
Description=Clears Completed Torrent

[Service]
User=ruxn
Group=ruxn
ExecStart=/opt/bin/cleartorrent.sh
Restart=no
