# systemd/movetocloud.service
[Unit]
Description=Cloud Uploader for Library Contents

[Service]
User=ruxn
Group=ruxn
ExecStart=/opt/bin/movetocloud.sh
Restart=no
