# systemd/torrentcleaner.service
[Unit]
Description=Torrent Cleaner

[Service]
User=ruxn
Group=ruxn
ExecStart=/opt/bin/torrentcleaner.sh
Restart=no
