# systemd/plexmediaserver.service.d/override.conf
[Unit]
Requires=mergerfs.service
After=network.target network-online.target mergerfs.service
